# postfix/supervisord.d/postfix.ini
[program:postfix]
process_name = postfix startup
command=postfix start
startsecs=0
autorestart=false

[program:postfix]
process_name = opendkim startup
command=opendkim
startsecs=0
autorestart=false