# postfix/supervisord.d/postfix.ini
[program:postfix]
process_name = postfix_startup
command=postfix start
startsecs=0
autorestart=false

[program:opendkim]
process_name = opendkim_startup
command=opendkim
startsecs=0
autorestart=false